# infra/systemd/hamalivpn-monitor.service
[Unit]
Description=HamaliVPN production health monitor
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hamalivpn-monitor
User=root
Group=root
UMask=0077
