# infra/systemd/hamali-offsite-backup.service
[Unit]
Description=HamaliVPN encrypted offsite backup
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/hamalivpn/infra/backup-offsite.sh
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=true
NoNewPrivileges=true
