# infra/systemd/hamali-external-uptime.service
[Unit]
Description=HamaliVPN external uptime check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/hamalivpn/uptime.env
ExecStart=/usr/bin/python3 /opt/hamali-uptime/external-uptime-monitor.py
User=hamali-uptime
Group=hamali-uptime
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hamali-uptime
