# infra/systemd/hamali-offsite-backup.timer
[Unit]
Description=Daily HamaliVPN encrypted offsite backup

[Timer]
OnCalendar=*-*-* 04:20:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=hamali-offsite-backup.service

[Install]
WantedBy=timers.target
